# home/systemd/.config/systemd/user/clean-backups.service
[Unit]
Description=Clean system/home backups

[Service]
Type=oneshot
ExecStart=/bin/bash %h/sync/system-data/clean.sh
