# home/systemd/.config/systemd/user/nvidia-settings.service
[Unit]
Description=Nvidia = pain

[Service]
Type=oneshot
ExecStart=/bin/bash %h/bin/nvidia.sh

[Install]
WantedBy=graphical.target
