# home/units/.config/systemd/user/backup.service
[Unit]
Description=Backup everything

[Service]
Type=oneshot
ExecStart=/bin/fish -c backup