# home/units/.config/systemd/user/backup.timer
[Unit]
Description=Backup everything daily

[Timer]
OnCalendar=daily
Persistent=true
#OnActiveSec=5m

[Install]
WantedBy=timers.target
