# home/systemd/.config/systemd/user/mpd.service.d/make-mpd-folders.conf
[Service]
ExecStartPre=/bin/sh -c "[ ! -d %h/.local/share/mpd ] && mkdir %h/.local/share/mpd || exit 0"
