# home/systemd/.config/systemd/user/mpd.service.d/make-mpd-folders.conf
[Service]
ExecStartPre=/bin/sh -c "[ ! -d %h/.local/share/mpd ] && mkdir -p %h/.local/share/mpd/playlists || exit 0"
